# debian/conf/cron.d/cosign-cleanup
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# Delete cosign files ( filter or daemon ) that are older than a day:
0 0 * * * root find /opt/cosign/db -type f -mtime +1 | xargs rm -f
# Done by monster:
#10 0 * * *  root  find /opt/cosign/db/daemon -type f -mtime +1 | xargs rm -f
